# init/supervisord.conf
[program:fireside]
directory=/home/fireside
command=fireside
autostart=true
autorestart=true
user=fireside
redirect_stderr=true
stdout_logfile=/home/fireside/fireside.log
stopwaitsecs=7